# melvin.service
[Unit]
Description=Melvin Humanoid Robot System
Documentation=https://github.com/melvin-robot
After=network.target
Wants=network.target

[Service]
Type=simple
User=melvin
Group=melvin
WorkingDirectory=/home/melvin/melvin
ExecStart=/usr/local/bin/melvin
ExecReload=/bin/kill -HUP $MAINPID
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal
SyslogIdentifier=melvin

# Security settings
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/var/melvin/logs /home/melvin/melvin/data
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true
RestrictRealtime=true
RestrictSUIDSGID=true

# Resource limits
LimitNOFILE=65536
LimitNPROC=4096
MemoryMax=2G
CPUQuota=200%

# Environment variables
Environment=LD_LIBRARY_PATH=/usr/local/lib
Environment=MELVIN_CONFIG_PATH=/etc/melvin
Environment=MELVIN_LOG_PATH=/var/melvin/logs
Environment=MELVIN_DATA_PATH=/home/melvin/melvin/data

[Install]
WantedBy=multi-user.target
