# archive/melvin.service
[Unit]
Description=Melvin Autonomous AI Runtime
After=network.target

[Service]
Type=simple
User=melvin
WorkingDirectory=/home/melvin/melvin_system
ExecStart=/home/melvin/melvin_system/melvin
Restart=always
RestartSec=10
StandardOutput=append:/home/melvin/melvin_system/melvin.log
StandardError=append:/home/melvin/melvin_system/melvin.log

[Install]
WantedBy=multi-user.target


